# scripts/scale_ingest.service
[Unit]
Description=Fitman scale ingestion — capture e.volve BLE reading and POST to API
After=network.target bluetooth.target

[Service]
Type=oneshot
# Update these paths to match your setup
WorkingDirectory=/home/YOUR_USER/Fitman
EnvironmentFile=/home/YOUR_USER/Fitman/.env
ExecStart=/home/YOUR_USER/Fitman/backend/.venv/bin/python scripts/scale_ingest.py
StandardOutput=journal
StandardError=journal
# Restart after each successful capture so the service is always ready
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
